# lasttask/renew.cron
0 4 * * * root /certbot.sh --renew >> /cronlog/renew.log 2>&1
